# backend/banhannah-backend.service
[Unit]
Description=BanHannah Educational Platform Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/banhannah/backend
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=banhannah-backend

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/banhannah/backend
ReadWritePaths=/home/ubuntu/banhannah/backend/storage
ReadWritePaths=/home/ubuntu/banhannah/backend/logs

# Resource limits for large uploads
LimitNOFILE=65536
MemoryMax=2G
# Increase timeout for large uploads
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target